# s6/amqp/run
#!/usr/bin/execlineb -P
python /app/run_amqp.py
