# s6/flask_n_ws/run
#!/usr/bin/execlineb -P
python /app/run.py
